# debian/filesystem-daemon.service
[Unit]
Description=Secure Filesystem Monitoring Daemon
Documentation=man:filesystem-daemon(8)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/filesystem-daemon
Restart=always
RestartSec=5
TimeoutStartSec=0
WorkingDirectory=/var/www/html

# Security settings
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
SystemCallFilter=~@debug
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Filesystem security
ReadWritePaths=/var/www/html
ReadOnlyPaths=/etc

[Install]
WantedBy=multi-user.target
