# debian/postinst
#!/bin/bash

set -e

# Crear directorio de trabajo si no existe
if [ ! -d "/var/www/html" ]; then
    mkdir -p /var/www/html
    chown root:root /var/www/html
    chmod 755 /var/www/html
fi

# Funci칩n para verificar si estamos en Docker
in_docker() {
    [ -f /.dockerenv ] || grep -q 'docker\|lxc' /proc/1/cgroup
}

# Funci칩n para verificar si systemd est치 disponible
has_systemd() {
    command -v systemctl >/dev/null 2>&1
}

# Solo intentar iniciar el servicio si:
# 1. No estamos en Docker
# 2. Systemd est치 disponible
if ! in_docker && has_systemd; then
    systemctl daemon-reload || true
    systemctl enable filesystem-daemon || true
    systemctl start filesystem-daemon || true
else
    echo "Skipping systemd service initialization (Docker environment or systemd not available)"
fi

exit 0
